# deploy/hiverelay/peerit-outbox-backup.service
[Unit]
Description=Back up Peerit OutboxLog atomic state
After=hiverelay.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/peerit-outbox-backup
Nice=10
IOSchedulingClass=idle
